# docker/supervisord/ppd/importPipelines.conf
[program:import-pipelines]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --queue=import-pipelines --sleep=3 --tries=3 --max-time=3600 --memory=1024
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/import-pipelines.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stopwaitsecs=3600
;MIN_REPLICAT=1
;MAX_REPLICAT=10
;QUEUE_THRESHOLD=2